# soppa/elasticsearch/config/elasticsearch_supervisor.conf
[program:elasticsearch]
user={{supervisor_user}}
directory={{path}}
command={{path}}bin/elasticsearch
autostart=true
autorestart=true
redirect_stderr=true
killasgroup=true

