# soppa/elasticsearch/config/elasticsearch_supervisor.conf
[program:elasticsearch]
user={{supervisor_user}}
directory={{project_root}}
command={{project_root}}bin/elasticsearch
autostart=true
autorestart=true
redirect_stderr=true
killasgroup=true

